# data/systemd/mobi.phosh.FileOpen.service.in
[Unit]
Description=A simple File Browser

[Service]
Type=dbus
BusName=@busname@
ExecStart=@bindir@/@bin@ --gapplication-service
KillMode=process
